# conf/bme280.conf
[program:bme280]
command=bash -c 'sleep 10 && /home/pi/python_apps/bme280-meshblu-py/bme280_publish.py'
numprocs=1
redirect_stderr=true 
stdout_logfile=/var/log/bme280.log
autorestart=true
user=pi
